# webapp/webapp/etc/cron.d/db-backup-cron
* * * * * /usr/local/bin/python /webapp/manage.py runcrons >/dev/null 2>&1
